# fastsearch-index-full.service
[Unit]
Description=FastSearch full consistency reindex
After=vps-fastsearch.service
Wants=vps-fastsearch.service

[Service]
Type=oneshot
ExecStart=%h/fastsearch/.venv/bin/python3 %h/fastsearch/scripts/fastsearch_indexer.py --mode full
Environment="FASTSEARCH_CONFIG=%h/.config/fastsearch/config.yaml"
Environment="FASTSEARCH_DB=%h/.local/share/fastsearch/fastsearch.db"
