# fastsearch.service
[Unit]
Description=FastSearch Memory Daemon
Documentation=https://github.com/your-username/fastsearch
After=network.target

[Service]
Type=simple
User=fastsearch
Group=fastsearch

# Main executable
ExecStart=/usr/local/bin/fastsearch daemon start
ExecStop=/usr/local/bin/fastsearch daemon stop
ExecReload=/usr/local/bin/fastsearch daemon reload

# Environment
Environment="FASTSEARCH_CONFIG=/etc/fastsearch/config.yaml"
Environment="FASTSEARCH_DB=/var/lib/fastsearch/fastsearch.db"

# Working directory
WorkingDirectory=/var/lib/fastsearch

# Restart policy
Restart=always
RestartSec=5

# Resource limits
MemoryMax=4G
MemoryHigh=3G

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/fastsearch /tmp

# Timeout
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
