# vps-fastsearch.service
[Unit]
Description=VPS-FastSearch Memory Daemon
Documentation=https://github.com/NarlySoftware/VPS-fastsearch
After=network.target

[Service]
Type=simple

# Main executable (venv at ~/venv)
ExecStart=%h/venv/bin/vps-fastsearch daemon start
ExecReload=%h/venv/bin/vps-fastsearch daemon reload

# Environment (install.sh puts config in ~/.config/fastsearch/)
Environment="FASTSEARCH_CONFIG=%h/.config/fastsearch/config.yaml"
Environment="FASTSEARCH_DB=%h/fastsearch/fastsearch.db"

# Working directory
WorkingDirectory=%h/fastsearch

# Restart policy
Restart=always
RestartSec=5

# Resource limits
MemoryMax=4G
MemoryHigh=3G

# Security hardening (user-service compatible)
NoNewPrivileges=true
PrivateTmp=false

# Timeout
TimeoutStartSec=120
TimeoutStopSec=30

[Install]
WantedBy=default.target
